# php/etc/supervisor/conf.d/nginx.ini
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autorestart = true
stderr_logfile = /tmp/nginx.err
stdout_logfile = /tmp/nginx.log
